# docker/supervisor/messenger-consume.ini
[program:messenger-consume]
directory=/srv/sharep
command=php ./bin/console messenger:consume async_priority_high async_priority_low --time-limit=60 --memory-limit=256MB  --limit=1
numprocs=1
autostart=true
autorestart=true
startsecs=5
startretries=10
process_name=%(program_name)s_%(process_num)02d
user=root
